# hr-advisor.service
[Unit]
Description=HR Policy Advisor Web App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hr-policy-advisor
ExecStart=/usr/bin/python3 /opt/hr-policy-advisor/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
